# root/etc/services.d/kodikontrol/run
#!/usr/bin/with-contenv bash

cd /app/kodikontrol || exit

exec \
	node ./bin/www
