# root/etc/services.d/kodikontrol/run
#!/usr/bin/with-contenv bash

cd /usr/lib/node_modules/KodiKontrol || exit

exec \
	node ./bin/www
