# root/etc/services.d/kodikontrol/run
#!/usr/bin/with-contenv bash

cd /usr/lib/node_modules/KodiKontrol/bin || exit

exec \
	s6-setuidgid abc \
	node www --home /config
